# scripts/com.swingpro.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  SwingPro health-check watchdog launchd agent.

  Probes /api/health/macro every 5 minutes. If the endpoint is
  unreachable, returns non-200, or shows scheduler stalled during
  market hours, kicks the main SwingPro launchd service to restart it.

  Install:
    cp scripts/com.swingpro.watchdog.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.swingpro.watchdog.plist
    launchctl list | grep swingpro.watchdog

  Logs: ~/Library/Logs/swingpro-watchdog.log
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.swingpro.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/arindamchowdhury/Development/Web Dev/Swing Stockpicker Prototype/scripts/healthcheck.sh</string>
    </array>

    <!-- Fire every 300 seconds (5 minutes). -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/arindamchowdhury/Library/Logs/swingpro-watchdog.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/arindamchowdhury/Library/Logs/swingpro-watchdog.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
